# GasTeminal/gas_station/build/gas_station.service
[Unit]

Description=gas_station
After=network.target

[Service]
Type=forking
ExecStart=/opt/GasStationPro/gas_station
ExecReload=/opt/GasStationPro/gas_station
WorkingDirectory=/home/vadosss63/.GasStationPro/
Restart=always

[Install]
WantedBy=graphical-session.target